# docker/services.d/graphql/finish
#!/bin/sh
#


echo 'graphql server died'
