# docker/services.d/graphql/run
#!/command/execlineb -P
#

s6-envdir /var/run/s6/container_environment
/bin/sh -f /app/scripts/graphql-server.sh
